# Modul6/systemd/pjk-tcp-server.service
[Unit]
Description=PJK TCP Server
After=network.target

[Service]
Type=simple
User=member
WorkingDirectory=/home/member/praktikum-pengantar-jaringan-komputer/Modul6/server
ExecStart=/usr/bin/python3 /home/member/praktikum-pengantar-jaringan-komputer/Modul6/server/tcp-server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
